# user_scripts/sddm/sddm/sddm-sync.path
[Unit]
Description=Watch for SDDM wallpaper sync requests

[Path]
PathModified=/tmp/sddm-sync-pending
Unit=sddm-sync.service

[Install]
WantedBy=multi-user.target
